# docker/python/build/3.8/development/supervisor/conf.d/ConsumerManager.conf
[program:ConsumerManager]
command=/app/venv/bin/python ConsumerManager.py
directory=/app/app/domain/Shopify
autostart=true
autorestart=true
# startretries=5
startsecs=0
user=root
numprocs=1
stderr_logfile=/var/log/supervisor/ConsumerManager.err.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/ConsumerManager.out.log
stdout_logfile_maxbytes=10MB
